# contrib/supervisord/template.ini
[program:PokeGoMapWorkerWRK]
command=/usr/bin/python DIRECTORYrunserver.py -ns -l "LOC" -st ST -sd SD -ld LD -a AUTH USER -p PASS
numprocs=1
directory=DIRECTORY
startsec=15
startretries=3
autorestart=true
stopwaitsecs=5
stdout_logfile=/tmp/supervisor_workerWRK.log
stderr_logfile=/tmp/supervisor_workerWRK.log
